# launchd/com.claude.memory-consolidation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.claude.memory-consolidation</string>
    <key>Comment</key>
    <string>Weekly memory consolidation for Claude Code self-learning system</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>-c</string>
      <string>/Users/psm2/Library/Mobile Documents/com~apple~CloudDocs/claude-setup/bin/memory-auto-consolidate.sh</string>
    </array>
    <!-- Run every Sunday at 3:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
      <key>Weekday</key>
      <integer>0</integer>
      <key>Hour</key>
      <integer>3</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/psm2/Library/Mobile Documents/com~apple~CloudDocs/claude-setup/memory/launchd-consolidation.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/psm2/Library/Mobile Documents/com~apple~CloudDocs/claude-setup/memory/launchd-consolidation-error.log</string>
    <key>RunAtLoad</key>
    <false/>
    <!-- Don't run if user is inactive -->
    <key>AbandonProcessGroup</key>
    <true/>
  </dict>
</plist>
