# launchd/com.claude.setup-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.claude.setup-sync</string>
    <key>Comment</key>
    <string>Auto-pull claude-setup from GitHub every 3 minutes</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>-c</string>
      <string>cd "$HOME/.claude-setup" &amp;&amp; git pull --ff-only --quiet origin master 2&gt;/dev/null || true</string>
    </array>
    <key>StartInterval</key>
    <integer>180</integer>
    <key>StandardOutPath</key>
    <string>/tmp/claude-setup-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/claude-setup-sync-error.log</string>
    <key>RunAtLoad</key>
    <true/>
  </dict>
</plist>
